# overlays/rpi-debian/etc/systemd/system/systemd-tmpfiles-setup-dev.service.d/walt.conf
[Service]
# Clear the LoadCredential directive or the unit fails in "walt image shell"
LoadCredential=
