# overlays/debian/etc/systemd/system/ptpd.service.d/walt.conf
[Unit]
# This should not run in "walt image shell"
ConditionVirtualization=!container
